# etc/cron.d/vmpr
*/15 * * * * root /opt/xensource/sm/runvmpr
